# board/air/rootfs-overlay/etc/inittab
# Air inittab
# Minimal init configuration

# Run init script at boot
::sysinit:/etc/init.d/rcS

# Reboot on Ctrl-Alt-Del
::ctrlaltdel:/sbin/reboot

# Shutdown properly
::shutdown:/bin/umount -a -r
